# mailcatcher/systemd/mailcatcher.service.jinja
[Unit]
Description=MailCatcher
After=network.service

[Service]
User=mailcatcher
Group=mailcatcher
ExecStart=/var/mailcatcher/.gem/ruby/{{ ruby_version }}/bin/mailcatcher --foreground --ip 0.0.0.0

[Install]
WantedBy=multi-user.target
